# quadlet/hermes-dashboard-pod.container
[Unit]
Description=Hermes Agent Dashboard (Pod)
After=network-online.target
Wants=network-online.target

[Container]
Image=quay.io/shanemcd/hermes-agent:vertex-55742
ContainerName=hermes-dashboard
Pod=hermes.pod

# Command
Exec=hermes dashboard --host 0.0.0.0 --no-open

# Environment
Environment=HOME=/opt/data
EnvironmentFile=%h/.config/hermes/env

# Volumes
Volume=hermes-data:/opt/data:Z
Volume=%h/.config/gcloud/application_default_credentials.json:/secrets/vertex/application_default_credentials.json:z,ro

# Resource limits
Memory=1G

[Service]
# Restart policy
Restart=always

[Install]
WantedBy=default.target
