# agents/openclaw/quadlet/nemoclaw-discord.container
[Unit]
Description=NemoClaw OpenClaw Discord Bot
After=network-online.target
Wants=network-online.target

[Container]
Image=localhost/nemoclaw-discord:latest
ContainerName=nemoclaw-discord

# Discord bot token
Secret=discord-bot-token,type=env,target=DISCORD_BOT_TOKEN

# OpenRouter API key
Secret=openrouter-api-key,type=env,target=OPENROUTER_API_KEY

# OpenRouter inference config
Environment=OPENROUTER_BASE_URL=https://openrouter.ai/api/v1
Environment=OPENCLAW_MODEL=anthropic/claude-sonnet-4-20250514

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=default.target
